# discord-bot.service
[Unit]
Description=Discord OpenRouter Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Discord-bot-new
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/Discord-bot-new/.venv/bin/python /home/pi/Discord-bot-new/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
